# roles/os_tuning/templates/profile.d/00-os-baseline.sh.j2
#!/bin/sh
# Managed by Ansible - do not edit manually

umask 0022
export LANG="en_US.UTF-8"
export LANGUAGE="en_US.UTF-8"
